# root/etc/services.d/aria2c/run
#!/usr/bin/with-contenv bash

exec $(eval echo "s6-setuidgid abc aria2c --conf-path=/data/config/custom.conf --enable-rpc=true --rpc-secret=${RPC_TOKEN} --rpc-listen-port=6800 -d /data/downloads/ ${OVERRIDE_OPTIONS}")
